# config/monitor/rc.local
/home/pi/bin/mqttmyip.sh wlan0 &
/home/pi/bin/mqttmyip.sh wlan1 &
/home/pi/bin/mqttmyip.sh eth0 &
exit 0
